# rpm/paxctld.service
[Unit]
Description=PaX flags maintenance daemon
After=syslog.target

[Service]
ExecStart=/sbin/paxctld

[Install]
WantedBy=multi-user.target
